# rootfs/etc/services.d/squid/run
#!/usr/bin/with-contenv bash
exec squid3 -f ${SQUID_CONF_DIR}/squid.conf -NYCd 1

